# keepalived/files/etc/hotplug.d/keepalived/553-other
#!/bin/sh

# shellcheck source=/dev/null
. /lib/functions/keepalived/hotplug.sh

# 记录所有环境变量，便于调试
env > /tmp/keepalived_hotplug_env.log

# 使用 ACTION 来判断状态
case "$ACTION" in
    NOTIFY_MASTER)
        echo "Keepalived state is MASTER, executing other_service_check.sh"
        logger "Keepalived state is MASTER, executing other_service_check.sh"
        /path/to/other_service_check.sh MASTER
        ;;
    NOTIFY_BACKUP)
        echo "Keepalived state is BACKUP, executing other_service_check.sh"
        logger "Keepalived state is BACKUP, executing other_service_check.sh"
        /path/to/other_service_check.sh BACKUP
        ;;
    *)
        echo "Unknown Keepalived action: $ACTION"
        logger "Unknown Keepalived action: $ACTION"
        ;;
esac
